# openstack-swift-proxy.upstart
description "OpenStack Object Storage (swift) - Proxy Server"

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [S016]

respawn

exec su -s /bin/sh -c "exec /usr/bin/swift-proxy-server /etc/swift/proxy-server.conf" swift
