# openstack-swift-container.upstart
description "OpenStack Object Storage (swift) - Container Server"

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [S016]

respawn

exec su -s /bin/sh -c "exec /usr/bin/swift-container-server /etc/swift/container-server.conf" swift
